# supervisord-mysql-new-db.conf
[program:mysqld]
command=/create-tematres-db.sh
numprocs=1
autostart=true
autorestart=true
